# rpi-monitor/setup/claude-monitor.service
[Unit]
Description=Claude Monitor Backend
Documentation=https://github.com/BorisLandoni/esp32-claude-token-monitor
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/claude-monitor/backend
Environment=PORT=8080
# Uncomment and set if using system Chromium (saves 300MB vs playwright install):
# Environment=CHROMIUM_PATH=/usr/bin/chromium-browser
ExecStart=/home/pi/claude-monitor/backend/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-monitor

[Install]
WantedBy=multi-user.target
